# config/includes.chroot/etc/update-motd.d/10-uname
#!/bin/sh
uname -snrvm
/usr/bin/neofetch